# deploy/tvlink.service
[Unit]
Description=TvLink Tavily API key pool
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tvlink
Group=tvlink
WorkingDirectory=/opt/tvlink
ExecStart=/opt/tvlink/tvlink -config /opt/tvlink/tvlink.toml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
